# docker/flask/gunicorn.conf
[program:gunicorn]
directory=/www/api
command=gunicorn -w 4 main:api
autostart=true
autorestart=true
stderr_logfile=/var/log/gunicorn/gunicorn.err.log
stdout_logfile=/var/log/gunicorn/gunicorn.out.log